# roles/node_exporter/templates/node_exporter.socket.j2
[Unit]
Description=Node Exporter

[Socket]
ListenStream={{ node_exporter_web_listen_port }}

[Install]
WantedBy=socket.target